# systemd/nyancat.socket
[Unit]
Description=badapple socket

[Socket]
ListenStream=23
Accept=true

[Install]
WantedBy=sockets.target
